# systemd/user/discord-monitor-scheme.path
[Unit]
Description=Sync discord and caelestia schemes (monitor)

[Path]
PathModified=%S/caelestia/scheme/current.txt
Unit=discord-monitor-scheme.service

[Install]
WantedBy=default.target
